# src/toilet.socket
[Unit]
Description=FAUST CTF toilet service

[Socket]
ListenStream=5743
Accept=true

[Install]
WantedBy=sockets.target
